# systemd/superdarn-sounder@.service
[Unit]
Description=superdarn-sounder passive SuperDARN monitor for radiod %i
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=notify
User=superdarn
Group=superdarn

WatchdogSec=180
NotifyAccess=main

EnvironmentFile=-/etc/sigmond/coordination.env
EnvironmentFile=-/etc/superdarn-sounder/env/%i.env

WorkingDirectory=/opt/git/sigmond/superdarn-sounder

ExecStartPre=+/bin/mkdir -p /var/lib/superdarn-sounder/%i
ExecStartPre=+/bin/mkdir -p /var/log/superdarn-sounder
ExecStartPre=+/usr/bin/chown -R superdarn:superdarn /var/lib/superdarn-sounder/%i
ExecStartPre=+/usr/bin/chown -R superdarn:superdarn /var/log/superdarn-sounder

# One systemd instance per signal source, each reporting under a unique
# reporter id (sigmond MULTI-INSTANCE-ARCHITECTURE.md §3): %i IS the reporter
# id after `smd instance migrate`, which is why the spool dir
# /var/lib/superdarn-sounder/%i (prepared above) is keyed on it.
# --instance %i picks up /etc/superdarn-sounder/%i.toml when present; otherwise
# the daemon falls back to the shared
# /etc/superdarn-sounder/superdarn-sounder-config.toml.  --radiod-id %i is the
# legacy fallback (radiod binding is config-driven via the [[radiod]] block).
# --config is intentionally omitted (explicit_path would pin the legacy path
# forever in resolve_config_path).
ExecStart=/opt/git/sigmond/superdarn-sounder/venv/bin/python3 -m superdarn_sounder.cli daemon \
    --instance %i \
    --radiod-id %i

TimeoutStartSec=60
Restart=always
RestartSec=5

StandardOutput=append:/var/log/superdarn-sounder/%i.log
StandardError=inherit
SyslogIdentifier=superdarn-sounder@%i

Nice=5
MemoryMax=512M
MemorySwapMax=0

NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/superdarn-sounder /var/log/superdarn-sounder
ReadOnlyPaths=/etc/superdarn-sounder
PrivateNetwork=false

CapabilityBoundingSet=CAP_NET_RAW CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_RAW CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
